# attachments/authorize_cluster_keys.cron
* * * * * root /home/rightlink/scripts/authorize_cluster_keys.sh
