# attachments/authorize_cluster_keys.cron
* * * * * root /home/rightscale/scripts/authorize_cluster_keys.sh
